# src/systemd/vigilant-proxy.service
[Unit]
Description=VIGILANT Transparent Proxy (mitmproxy)
After=network.target dnsmasq.service
Wants=dnsmasq.service

[Service]
Type=simple
User=vigilant_admin
WorkingDirectory=/home/vigilant_admin/vigilant
ExecStart=/home/vigilant_admin/vigilant/venv/bin/mitmdump \
  --mode transparent \
  --showhost \
  --listen-host 0.0.0.0 \
  --listen-port 8080 \
  --set block_global=false \
  -s /home/vigilant_admin/vigilant/addons/vigilant_addon.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/vigilant_admin/vigilant/logs/proxy.log
StandardError=append:/home/vigilant_admin/vigilant/logs/proxy.log

[Install]
WantedBy=multi-user.target
